# PMR_DC_3/codetar/etc/PMR_DC.cron
MAILTO=supreet.singh@guavus.com

# Run the 15 minute job to get all node statuses in DC
*/5 * * * *    root    /data/scripts/PMR/insertdcvalue/bin/fiveMinute.sh

# Run the 15 minute job to get all node statuses in DC 
*/15 * * * *	root	/data/scripts/PMR/insertdcvalue/bin/fifteenMinute.sh

# Run the hourly job to get CPU/Mem Utilization from all DC Nodes
59 * * * *	root	/data/scripts/PMR/insertdcvalue/bin/hourly.sh

# Run the Sync Job every day at 30 minutes past the 00 hour, determine OS build version and patches
30 00 * * *      root    /data/scripts/PMR/insertdcvalue/bin/daily.sh

# Run the Sync Job every hour at 35 minutes past the hour, push data from DC to SAP
35 * * * *      root    /data/scripts/PMR/insertdcvalue/bin/SyncDC_PMRHosts.sh

# Run the Sync Job every day at 30 minutes past the 23 hour, download code from SAP to DC
30 23 * * *      root    /data/scripts/PMR/insertdcvalue/bin/SyncSAP_DCHosts.sh

# Clean up older than 7 days files
55 23 * * *     root    /data/scripts/PMR/insertdcvalue/bin/Cleanup.sh

# Sanitize the site architecture.
00 01 * * *      root	/data/scripts/PMR/insertdcvalue/bin/prepare_site.sh
